# docker/8_angular_flask_go/config/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:go]
command=/var/www/goserver/main